# systemd/user/udiskie.service
[Unit]
Description=UDiskie Automounter

[Service]
ExecStart=/usr/bin/udiskie

[Install]
WantedBy=vodik.target
