# systemd/user/vodik.target
[Unit]
Description=Stuff for Vodik
After=xorg.target

[Install]
Alias=default.target
